# src/scheduling/systemd/timers/system/package-cache/package-cache.service
[Unit]
Description=Clean pacman package cache, keep last 3 versions

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -r
